# tools/perception/yolo_follow/oneday-yolo-follow.service
[Unit]
Description=Project4 YOLO Person Follow
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oneday
Group=oneday
WorkingDirectory=/home/oneday/project4/tools/perception/yolo_follow
ExecStart=/home/oneday/project4/tools/perception/yolo_follow/run_follow_target.sh
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
